# LaunchAgents/com.example.jarvis.ui.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.jarvis.ui</string>

    <!--
      Control Room dashboard (jarvis-ui-server.py) en daemon KeepAlive.
      Pattern /bin/zsh -c "<wrapper>" — TCC-safe, hérite des grants.
      Le wrapper fait `nohup python &` mais avec --no-open pour ne pas spawn de browser.
      Sans le KeepAlive, un crash silencieux laissait le dashboard servir du stale jusqu'à
      restart manuel (cf. régression observée le 2026-05-09).
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>exec /usr/bin/python3 __HOME__/.local/bin/jarvis-ui-server.py --port 7474</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>__HOME__/.local/var/log/jarvis-ui.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.local/var/log/jarvis-ui.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PATH</key>
        <string>__HOME__/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
</dict>
</plist>
