# LaunchAgents/com.example.jarvis.vault-index-build.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.jarvis.vault-index-build</string>

    <!--
      Régénère les 2 index markdown auto-loadés (HOT) :
        - _vault-index.md (carte du vault)       via jarvis-vault-index-build
        - Memory/auto/_index.md (skills auto)    via jarvis-skills-bundle
      Distinct de com.example.jarvis.vault-index (embeddings sqlite-vec).
      Builders Python purs, idempotents, 0 invocation Anthropic.
      Remplace le déclencheur perdu : jarvis-brief.sh (supprimé le 2026-05-14).
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>__HOME__/.local/bin/jarvis-vault-index-build; __HOME__/.local/bin/jarvis-skills-bundle</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/.local/var/log/jarvis-vault-index-build.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.local/var/log/jarvis-vault-index-build.stderr.log</string>

    <key>WorkingDirectory</key>
    <string>__HOME__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PATH</key>
        <string>__HOME__/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
